# etc/profile.d/aws.sh
export EC2_HOME=/usr/local/ec2-api-tools
export S3_HOME=/usr/local/s3sync
export PATH=$PATH:$EC2_HOME/bin:$S3_HOME:
