# etc/profile.d/git.sh
alias ga='git add'
alias gb='git branch'
alias gc='git commit'
alias gd='git diff'
alias go='git checkout'
alias gs='git status'
alias gf='git fetch'
alias gfr='git fetch && git rebase'
alias gl='git l'
alias gm='git merge'
alias gp='git push'
alias gr='git rebase'
alias gt='git reset'
alias gu='git pull'
alias gst='git stash'
alias gsp='git stash pop'
