# etc/profile.d/misc.sh
alias cls='clear'
alias hgrep='history | grep'
alias rgrep='rake -T | grep'
alias h='history | tail -n 24'
alias mysql='mysql -p'
